# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:cron]
command=/usr/sbin/crond -f -l 8